# deploy/goliath-bot.service
[Unit]
Description=GOLIATH Telegram Bot — Autonomous Solar Portfolio Management
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=goliath
Group=goliath
WorkingDirectory=/opt/goliath/telegram-bot
ExecStart=/opt/goliath/venv/bin/python -m bot.main
ExecStop=/bin/kill -SIGTERM $MAINPID

# Environment
EnvironmentFile=/opt/goliath/.env

# Auto-restart on crash
Restart=always
RestartSec=10

# Logging (goes to journald; view with: journalctl -u goliath-bot -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=goliath-bot

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/opt/goliath /home/goliath
PrivateTmp=yes

# Resource limits
MemoryMax=1.5G
CPUQuota=80%

# Note: No WatchdogSec — the bot doesn't implement sd_notify heartbeats.
# Health monitoring is handled by the healthcheck.sh cron job instead.

[Install]
WantedBy=multi-user.target
